# fs/etc/services.d/grafana/run
#!/usr/bin/with-contenv ash

cd /usr/share/grafana

exec runas ash /run.sh